# deploy/docker/crontab
* * * * * /bin/bash -c '/usr/src/qatrackplus/deploy/docker/cron_backup.sh'
