# deploy/docker/crontab
0 */6 * * * /usr/src/qatrackplus/deploy/docker/cron_backup.sh >> /usr/src/qatrackplus/deploy/docker/user-data/cron.log 2>&1
